# Auto_Config/base/crontab
## Atualiza o sistema as 20:00, segunda, quarta e sexta.
00 20 * * 1,3,5 apt update ; apt upgrade
