# Auto_Config/base/crontab
## Atualiza o sistema as 20:00, segunda, quarta e sexta.
00 20 * * 01,03,05 apt update > /tmp/crontab.txt ; apt upgrade -y >> /tmp/crontab.txt && date >> /tmp/crontab.txt